# skills/model-failover-guard/openclaw-model-failover.service
[Unit]
Description=OpenClaw Model Failover Guard
After=network.target

[Service]
Type=simple
User=%u
WorkingDirectory=%h
ExecStart=/usr/bin/python3 %h/.openclaw/workspace/skills/model-failover-guard/skills/model-failover-guard/scripts/failover.py loop
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target